# deploy/swinglogger.service
[Unit]
Description=SwingLogger Flask app (Baseline Tech)
After=network.target

[Service]
Type=simple
# EDIT THESE TO MATCH YOUR VPS:
User=__USER__
Group=__USER__
WorkingDirectory=__HOME__/swinglogger
Environment="PATH=__HOME__/swinglogger/.venv/bin"
ExecStart=__HOME__/swinglogger/.venv/bin/gunicorn \
    --workers 2 \
    --bind 127.0.0.1:8001 \
    --timeout 60 \
    --access-logfile __HOME__/swinglogger/access.log \
    --error-logfile __HOME__/swinglogger/error.log \
    app:app

Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
